# profile.d/services.sh
export PLUSHU_SERVICES_DIR=$PLUSHU_DIR/services
